# home/.config/systemd/user/dynamite.service
[Unit]
Description=Dynamically change wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dynamite-update
