# home/.config/systemd/user/dynamite.timer
[Unit]
Description=Dynamite Interval

[Timer]
OnUnitActiveSec=5min
OnBootSec=5min

[Install]
WantedBy=timers.target
